# setup/aws/detrack.service
[Unit]
Description=Detrack GPS tracking server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/detrack
Environment=CONFIG_USE_ENVIRONMENT_VARIABLES=true
EnvironmentFile=-/opt/detrack/conf/detrack.env
ExecStart=/usr/bin/java -Xms256m -Xmx768m -jar detrack-server.jar conf/traccar.xml
SyslogIdentifier=detrack
SuccessExitStatus=143
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
